# templates/default/dela.service.erb
[Unit]
Description = Dela server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.dela.user %>

PIDFile = /tmp/dela.pid
ExecStart = <%= node.dela.home %>/bin/start.sh
ExecStop = <%= node.dela.home %>/bin/stop.sh
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
